# service/com.nrhof.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nrhof.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nrhof/dev/nrhof-bot/.venv312/bin/python</string>
        <string>/Users/nrhof/dev/nrhof-bot/app.py</string>
        <string>--fullscreen</string>
        <string>--display</string>
        <string>1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nrhof/dev/nrhof-bot</string>

    <!-- RunAtLoad disabled - manually start with 'make mac-start' -->
    <key>RunAtLoad</key>
    <false/>

    <!-- KeepAlive disabled - won't restart on exit -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/nrhof/dev/nrhof-bot/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nrhof/dev/nrhof-bot/logs/bot.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
