# service/nrhof-matrix-kiosk.service
[Unit]
Description=NRHOF Matrix Kiosk
After=network.target sound.target graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/nrhof-matrix-kiosk
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/home/pi/nrhof-matrix-kiosk/.venv312/bin/python /home/pi/nrhof-matrix-kiosk/app.py --fullscreen
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
